# socatlord/systemd/socatlord.service
[Unit]
Description=SOCAT manager

[Service]
Type=oneshot
ExecStart=%s run
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
